# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

; This is where you run individual Tornado instances.
; We run four; one per processor core.
; In development, we ran as many as four per core with no issues.
; If you're looking to minimize cpu load, run fewer processes.
; BTW, Tornado processes are single threaded.
; To take advantage of multiple cores, you'll need multiple processes.

[program:kinit]
command=/home/superset/kinit.sh
stderr_logfile = /var/log/supervisord/kinit.err.log 
stdout_logfile = /var/log/supervisord/kinit.out.log 
autorestart=true

[program:main]
command=/usr/local/bin/gunicorn --timeout ${GUNICORN_TIMEOUT:-120} --workers ${GUNICORN_WORKERS:-10} superset:app
stderr_logfile = /var/log/supervisord/gunicorn.err.log 
stdout_logfile = /var/log/supervisord/gunicorn.out.log 
autorestart=true

[program:flower]
command=/usr/local/bin/celery flower --app=superset.sql_lab:celery_app
stderr_logfile = /var/log/supervisord/flower.err.log 
stdout_logfile = /var/log/supervisord/flower.out.log 
autorestart=true

[program:worker1]
command=/usr/local/bin/celery worker --app=superset.sql_lab:celery_app --pool=gevent -Ofair -nworker1
stderr_logfile = /var/log/supervisord/worker1.err.log 
stdout_logfile = /var/log/supervisord/worker1.out.log 
autorestart=true

[program:worker2]
command=/usr/local/bin/celery worker --app=superset.sql_lab:celery_app --pool=gevent -Ofair -nworker2
stderr_logfile = /var/log/supervisord/worker2.err.log 
stdout_logfile = /var/log/supervisord/worker2.out.log 
autorestart=true

[program:worker3]
command=/usr/local/bin/celery worker --app=superset.sql_lab:celery_app --pool=gevent -Ofair -nworker3
stderr_logfile = /var/log/supervisord/worker3.err.log 
stdout_logfile = /var/log/supervisord/worker3.out.log 
autorestart=true
